# systemd/kudos-dashboard.service
[Unit]
Description=Kudos Shiny Dashboard
After=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/kudos
EnvironmentFile=/opt/kudos/.env
ExecStart=/opt/kudos/.venv/bin/gunicorn dash_app:server --bind 0.0.0.0:8050
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
