# systemd/kudos-backfill.service
[Unit]
Description=Kudos backfill embeddings and clusters
After=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/kudos
EnvironmentFile=/opt/kudos/.env
ExecStart=/opt/kudos/.venv/bin/python -m cron.backfill
ExecStartPost=/usr/bin/systemctl restart kudos-dashboard.service
