# systemd/kudos-accounting.service
[Unit]
Description=Kudos monthly accounting report
After=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/kudos
EnvironmentFile=/opt/kudos/.env
ExecStart=/opt/kudos/.venv/bin/python -m cron.accounting
